# certbot/templates/certbot.service.j2
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet --webroot -w /var/www/letsencrypt --deploy-hook "systemctl reload nginx"
